# deployment/systemd/powerpipeline-spp-ingest.service.template
[Unit]
Description=PowerPipeline SPP Load Ingestion (mtlf-vs-actual)
Documentation=file:/home/deboerja/git/github/PowerPipeline/docs/OPERATIONS.md

[Service]
Type=oneshot
Environment=POWERPIPELINE_RUNTIME_ROOT=/srv/apps/powerpipeline
ExecStart=/home/deboerja/git/github/PowerPipeline/scripts/run_spp_ingest.sh
Nice=5
TimeoutStartSec=120

# This is a template only. Do not install or enable until explicitly
# approved -- see docs/MVP_AUTONOMY_PROFILE.md agent-imposed checkpoints.
