# implementation/EXECUTION_HISTORY.yaml
version: 1
execution_history:
  - timestamp: "2026-07-16T22:37:00Z"
    event: repository_located
    detail: "Confirmed deboerja/PowerPipeline on GitHub via api.github.com; public, empty, default branch main"

  - timestamp: "2026-07-16T22:40:00Z"
    event: repository_cloned
    detail: "Cloned to /home/deboerja/git/github/PowerPipeline; push credential confirmed present via git credential fill"

  - timestamp: "2026-07-16T22:45:00Z"
    event: reuse_discovery_completed
    detail: "Explore-agent inventory of Enphase, weather, ResidentAI wrapper, scheduling, secrets, and runtime-layout conventions; written to docs/EXISTING_COMPONENT_REUSE.md"

  - timestamp: "2026-07-16T22:50:00Z"
    event: spp_source_validation_completed
    detail: "Confirmed mtlf-vs-actual, stlf-vs-actual, generation-mix-historical reachable and downloadable without auth; renewable forecast-vs-actual dataset not located after ~25 probe attempts, recorded as bl-001"

  - timestamp: "2026-07-16T22:54:00Z"
    event: feature_branch_created
    detail: "feature/mvp-bootstrap-v1"

  - timestamp: "2026-07-16T23:00:00Z"
    event: control_plane_docs_written
    detail: "NORTH_STAR, MVP_CHARTER, MVP_AUTONOMY_PROFILE (md+yaml), SOURCE_REGISTRY, DECISION_LOG, EXISTING_COMPONENT_REUSE, ARCHITECTURE, SECURITY_AND_AUTHORITY, DATA_LINEAGE, OPERATIONS, FAILURE_SCENARIOS, INTERVIEW_DEMO, FUTURE_ARCHITECTURE, and implementation/*.yaml files written"

  - timestamp: "2026-07-16T23:07:00Z"
    event: dev_environment_set_up
    detail: "python venv (.venv) created, pyproject.toml deps (duckdb, pandas, pyarrow, pandera, httpx, pydantic, pytz) installed"

  - timestamp: "2026-07-16T23:15:00Z"
    event: spp_load_vertical_slice_implemented
    detail: >
      src/powerpipeline/contracts/spp_mtlf.py (Pandera contract),
      src/powerpipeline/ingest/spp_load.py (fetch/land_raw/normalize),
      src/powerpipeline/pipeline.py (orchestration + DuckDB upsert +
      pipeline_runs/data_quality_results/source_freshness), sql/schema.sql,
      sql/views.sql. 13 tests written and passing across
      tests/test_spp_mtlf_contract.py, tests/test_spp_load_ingest.py,
      tests/test_pipeline_curated.py, tests/test_missing_interval_detection.py.
      ruff check clean.

  - timestamp: "2026-07-16T23:18:00Z"
    event: live_smoke_test_passed
    detail: >
      Ran src/powerpipeline/pipeline.py against the real, live SPP endpoint
      (not a fixture) for 2026-07-16's OP-MTLF file: 386/386 records
      accepted, 0 rejected. Confirms the ingestion code works against
      production data, not just fixtures.

  - timestamp: "2026-07-16T23:19:00Z"
    event: feature_branch_pushed
    detail: "git push -u origin feature/mvp-bootstrap-v1 (commit fbe24b7)"

  - timestamp: "2026-07-16T23:35:00Z"
    event: household_bridges_implemented
    detail: >
      Sanitized fixtures created for Enphase daily-summary, weather
      daily-actual, and household solar forecast (synthetic values, real
      schema, deliberately not real household numbers per the
      no-unsanitized-telemetry boundary). Implemented
      src/powerpipeline/contracts/{enphase_summary,weather_actual,
      household_solar_forecast}.py and
      src/powerpipeline/ingest/{enphase_bridge,weather_bridge,
      household_solar_forecast_bridge}.py, plus three new pipeline.py
      orchestration functions. 9 new tests across
      tests/test_household_bridges.py and tests/test_household_contracts.py,
      all passing; ruff clean.

  - timestamp: "2026-07-16T23:42:00Z"
    event: residentai_tool_implemented
    detail: >
      src/powerpipeline/residentai_tool.py with 8 named operations
      (read_household_solar_history, read_household_solar_forecast,
      read_household_forecast_accuracy, read_spp_load_conditions,
      read_spp_load_forecast_accuracy, read_spp_renewable_conditions
      [honestly reports not_available], read_pipeline_health,
      read_data_quality_summary). 11 tests in test_residentai_tool.py,
      including structural verification that DuckDB's read_only=True
      connection rejects DELETE/DROP/INSERT and that no query/sql method
      exists on the class at all. 31 tests total, all passing; ruff clean.
      Not registered with OpenWebUI/ResidentAI -- that stays behind the
      deployment checkpoint.

  - timestamp: "2026-07-16T23:47:00Z"
    event: scheduling_templates_and_retry_backfill_implemented
    detail: >
      6 systemd unit templates under deployment/systemd/ (spp-ingest,
      household-bridge, quality-check -- each service+timer), all pass
      systemd-analyze verify. Corrected docs/OPERATIONS.md, which had
      described a nonexistent "recompute forecast-vs-actual views" timer --
      those are plain SQL views, not materialized tables, so there's nothing
      to recompute on a schedule. Added retry-with-backoff to
      ingest/spp_load.py::fetch_raw_csv and pipeline.bounded_backfill().
      Writing the failure-recovery tests surfaced and fixed a real bug: a
      fetch failure was raising before run_spp_load_ingest's try block, so
      it would never have been recorded in pipeline_runs -- restructured so
      failures are always recorded. 4 new tests in test_failure_recovery.py,
      2 in test_quality_sweep.py.

  - timestamp: "2026-07-16T23:50:00Z"
    event: interview_demo_built_and_executed
    detail: >
      src/powerpipeline/demo.py implementing all 20 steps from
      docs/INTERVIEW_DEMO.md. Actually run (python -m powerpipeline.demo
      run-all), not just written -- caught a second real issue during that
      run: the missing-interval-detection step was reusing fixture
      timestamps already loaded by earlier steps, so the "gap" it introduced
      was silently backfilled by prior upserts before the completeness check
      ran, which would have shown misleading demo output. Fixed by isolating
      that one demonstration in a fresh runtime root, matching the pattern
      already used by tests/test_missing_interval_detection.py. Reran and
      confirmed the gap is now correctly flagged. Generated
      reports/interview_demo_report.html. 37 tests total, all passing; ruff
      clean.

  - timestamp: "2026-07-16T23:56:00Z"
    event: user_approved_srv_deployment
    detail: "User: 'go ahead, deploy to /srv on Odin' -- explicit go-ahead for the deployment checkpoint (does not cover ResidentAI restart or interface registration, which remain separate checkpoints)."

  - timestamp: "2026-07-16T23:58:00Z"
    event: runtime_root_revised_mid_deployment
    detail: >
      mkdir /srv/powerpipeline failed with permission denied (/srv is
      root:root, mode 755; sudo requires a password this session doesn't
      have). Pivoted to /srv/apps/powerpipeline/ (owned by deboerja,
      matches every other homelab service's convention exactly, already
      scoped as the fallback in docs/EXISTING_COMPONENT_REUSE.md before
      deployment even started). Updated all docs, .env.example, and systemd
      templates to reference the new path; re-verified with
      systemd-analyze verify (clean). Recorded in docs/DECISION_LOG.md.

  - timestamp: "2026-07-17T00:00:00Z"
    event: srv_directory_and_env_created
    detail: "/srv/apps/powerpipeline/{raw,quarantine,normalized,database,metadata,logs,reports,backups,state} created; ~/.config/residentai/powerpipeline.env written (mode 600, path config only, no secrets)."

  - timestamp: "2026-07-17T00:01:00Z"
    event: systemd_units_installed_and_enabled
    detail: "6 unit files copied to ~/.config/systemd/user/ (template disclaimer replaced with an approval note); daemon-reload; enabled + started all 3 timers."

  - timestamp: "2026-07-17T00:02:57Z"
    event: spp_ingest_auto_run_correctly_failed_closed
    detail: >
      The spp-ingest service (manually triggered to validate) requested
      2026-07-17's file, which SPP had not yet published (confirmed via
      direct API check -- the folder was empty). Correctly raised, and the
      failure was correctly recorded as a 'failed' pipeline_runs row with
      source_freshness.status='failing' -- not hidden, not fabricated.

  - timestamp: "2026-07-17T00:05:27Z"
    event: household_bridge_auto_run_succeeded_on_real_schedule
    detail: >
      powerpipeline-household-bridge.timer fired on its own real schedule
      (not manually triggered) and succeeded against real production data:
      420 household solar rows, 421 weather rows, 7 solar-forecast rows
      landed in the deployed DuckDB at /srv/apps/powerpipeline/database/
      powerpipeline.duckdb.

  - timestamp: "2026-07-17T00:06:13Z"
    event: spp_ingest_manually_validated_against_real_existing_file
    detail: >
      Ran pipeline.run_spp_load_ingest against OP-MTLF-202607161800.csv
      (confirmed to actually exist via direct API check) on the deployed
      system: 386/386 records accepted, 0 rejected. Proves the SPP path
      itself works correctly against production; the earlier failure was
      purely an upstream publication-timing condition, not a bug.

  - timestamp: "2026-07-17T00:08:00Z"
    event: full_residentai_tool_validation_against_deployed_database
    detail: >
      Called read_household_solar_history, read_household_solar_forecast,
      read_spp_load_forecast_accuracy, read_pipeline_health, and
      read_data_quality_summary directly against the deployed database --
      all returned real, correct data. Re-verified the write-rejection
      guarantee against the deployed database specifically (a DELETE via
      the tool's own connection raised DuckDB's read-only error).
      Acceptance criteria upgraded from met_in_dev to met accordingly (see
      implementation/MVP_ACCEPTANCE_CRITERIA.yaml).
